# resources/systemd/reflector.timer
[Unit]
Description= Run reflector daily for https rated mirrors for speed

[Timer]
OnCalendar=*-*-* 22:00:00
AccuracySec=1s
Persistent=true
